# scripts/deploy/bxc.service
[Unit]
Description=Bxc Browser Automation API Server
After=network.target

[Service]
# Environment setup
# adjust WorkingDirectory/ExecStart to your deploy path
Type=simple
User=ubuntu
WorkingDirectory=/opt/bxc
ExecStart=/usr/bin/bun run src/cli/serve.ts
Restart=always
RestartSec=5

# Resources & Security
MemoryMax=2G
CPUWeight=100
StandardOutput=append:/var/log/bxc/api.log
StandardError=append:/var/log/bxc/error.log

# Ensure /var/log/bxc exists and is writable
# Environment=PORT=3000
# Environment=BXC_DB_PATH=./data/bxc.sqlite

[Install]
WantedBy=multi-user.target
